# meta-tibbo/recipes-multimedia/pulseaudio/pulseaudio/etc/systemd/system/pulseaudio.service
[Unit]
Description=Sound Service

[Service]
#ExecStart=/usr/bin/pulseaudio --system --disallow-exit -vvvvvvvvvvvvvvv
ExecStart=/usr/bin/pulseaudio --system --disallow-exit
Restart=on-failure

[Install]
WantedBy=multi-user.target
